# infra/docker/gui/supervisord.conf
; SPDX-License-Identifier: MIT
; Process supervisor for sevn gateway GUI image (Xvfb + noVNC + uvicorn).

[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:desktop]
command=/opt/sevn/infra/docker/gui/start-desktop.sh
user=sevnoperator
environment=DISPLAY=":99",HOME="/home/sevnoperator"
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:gateway]
command=/opt/sevn/.venv/bin/uvicorn sevn.gateway.http_server:create_app --factory --host 0.0.0.0 --port 3001
directory=/operator/workspace
user=sevnoperator
environment=HOME="/home/sevnoperator",SEVN_HOME="/operator",PATH="/opt/sevn/.venv/bin:/usr/local/bin:/usr/bin:/bin",DISPLAY=":99"
autorestart=true
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
